# deploy/production/backup.cron.example
0 3 * * * cd /opt/sellgram/deploy/production && BACKUP_DIR=/opt/sellgram/backups ./backup-db.sh >> /var/log/sellgram-backup.log 2>&1
