# deploy/production/monitor-health.cron.example
*/5 * * * * cd /opt/sellgram/deploy/production && ./monitor-health.sh >> /var/log/sellgram-health.log 2>&1
